# file_root/desktop/fcitx/files/etc/profile.d/exz-fcitx.sh
export XMODIFIERS="@im=fcitx"
export XIM=fcitx
export XIM_PROGRAM=fcitx
export GTK_IM_MODULE="XIM"
export QT_IM_MODULE="XIM"
